# systemd/munin-node-daemon.service
[Unit]

[Service]
ExecStart=/usr/local/sbin/munin-node-c -a

PrivateTmp=yes
PrivateNetwork=no
ProtectHome=yes

[Install]
WantedBy=multi-user.target
