# systemd/munin-node-c.socket
[Socket]
ListenStream=4949
Accept=true

[Install]
WantedBy=sockets.target
